# s6-overlay/s6-rc.d/sssd/run
#!/command/execlineb -P
/usr/sbin/sssd -i -f
