# s6-overlay/s6-rc.d/wait/run
#!/command/execlineb -P
echo "wait executing"
