# s6-overlay/s6-rc.d/populateusers/run
#!/command/execlineb -P
echo "populateusers executing"
